# deploy/bit-code-lab.service
[Unit]
Description=BIT Code Lab (FastAPI + built frontend)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/bit-code-lab/backend
ExecStart=/home/ubuntu/bit-code-lab/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
